# backend/configs/etc/systemd/system/fcgiwrap.service.d/override.conf
[Service]
Environment="DAEMON_OPTS=-f -c 2"
ExecStartPre=+/usr/lib/wb-homeui-backend/fcgiwrap_prestart.sh
ExecStart=
ExecStart=/usr/sbin/fcgiwrap $DAEMON_OPTS
