# backend/configs/etc/systemd/system/fcgiwrap.service.d/override.conf
[Service]
Environment="DAEMON_OPTS=-f -c 2"
ExecStartPre=+/usr/sbin/usermod -aG mosquitto www-data
ExecStart=
ExecStart=/usr/sbin/fcgiwrap $DAEMON_OPTS
