# backend/configs/usr/lib/systemd/system/fcgiwrap.service.d/10-mosquitto-group.conf
[Service]
ExecStartPre=+/usr/sbin/usermod -aG mosquitto www-data
