# backend/configs/usr/lib/systemd/system/fcgiwrap.service.d/10-prefork-2-procs.conf
[Service]
Environment="DAEMON_OPTS=-f -c 2"
ExecStart=
ExecStart=/usr/sbin/fcgiwrap $DAEMON_OPTS
