# backend/configs/usr/lib/systemd/system/mosquitto.service.d/20-fix-homeui-acl-perms.conf
[Service]
ExecStartPre=+/bin/chmod 0700 /usr/share/wb-mqtt-homeui/mosquitto/acl/user.conf
ExecStartPre=+/bin/chmod 0700 /usr/share/wb-mqtt-homeui/mosquitto/acl/operator.conf
ExecStartPre=+/bin/chown mosquitto:mosquitto /usr/share/wb-mqtt-homeui/mosquitto/acl/user.conf
ExecStartPre=+/bin/chown mosquitto:mosquitto /usr/share/wb-mqtt-homeui/mosquitto/acl/operator.conf
